# ansible/alpine/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=init

[program:ssh]
command=rc-service sshd start

[program:bash]
command=bash
