# ansible/alpine/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/bin/sh -c "/init.sh && supervisorctl start sh"
autorestart=false

[program:sh]
command=sh
autostart=false
