# ansible/ubuntu/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=service ssh start

[program:bash]
command=bash
